# files/etc/cron.d/update-monkeysphere-auth
*/5 * * * * root /usr/sbin/monkeysphere-authentication update-users
